# config/movie.conf
[program:movie]
command=/home/john/venv/bin/gunicorn djangoProject.wsgi:application -c /home/john/djangoProject/config/gunicorn.conf.py
directory=/home/john/djangoProject
user=john
autorestart=true
redirect_stderr=true
stdout_logfile = /home/john/djangoProject/logs/debug.log
